# templates/usr/lib/systemd/system/ont-reboot.service.tmpl
[Unit]
Description=ONT Reboot Service
ConditionPathExists={{ home_ }}/ont-reboot/ont-reboot.py
After=network.target

[Service]
ExecStart=python3 {{ home_ }}/ont-reboot/ont-reboot.py

[Install]
WantedBy=multi-user.target
